# systemd/bee_cashout.service
[Unit]
Description=Bee Cashout Service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/dailycash.sh

[Install]
WantedBy=default.target
