# systemd/bee_cashout.timer
[Unit]
Description=Bee Cashout Daily Routine

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
